# git_crontab.cron
0,1 * * * * sh /usr/local/yele-server/git_crontab.sh
