# git_crontab.cron
*/2 * * * * sh /usr/local/yele-server/git_crontab.sh
